# configs/systemd/update.timer
[Unit]
Description=Run VPS Updates Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
